# recipe_proxy/rootfs/etc/s6-overlay/s6-rc.d/init-nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Init nginx by processing templates
# ==============================================================================

declare entry
declare frontend_server
declare backend_server
declare proxy_pass_host

entry=$(bashio::addon.ingress_entry)
frontend_server=$(bashio::config 'frontend_server')
backend_server=$(bashio::config 'backend_server')
# The caret stops bashio:var.json from passing the value as a string
proxy_pass_host=^$(bashio::config 'proxy_pass_host')
proxy_pass_real_ip=^$(bashio::config 'proxy_pass_real_ip')

bashio::log.info "Ingress entry: ${entry}"
bashio::log.info "Frontend server: ${frontend_server}"
bashio::log.info "Backend server: ${backend_server}"
bashio::log.info "Proxy pass host: ${proxy_pass_host}"
bashio::log.info "Proxy pass real IP: ${proxy_pass_real_ip}"

bashio::var.json \
    entry "${entry}" \
    frontend_server "${frontend_server}" \
    backend_server "${backend_server}" \
    proxy_pass_host "${proxy_pass_host}" \
    proxy_pass_real_ip "${proxy_pass_real_ip}" \
    | tempio \
        -template /etc/nginx/templates/ingress.gtpl \
        -out /etc/nginx/servers/ingress.conf

bashio::log.info "Generated nginx configuration:"
cat /etc/nginx/servers/ingress.conf
